# roles/mesos/files/mesos-slave.conf
description "mesos slave"

start on started docker
stop on stopping docker

respawn
respawn limit 10 10

exec /usr/bin/mesos-init-wrapper slave
