# roles/monitored/templates/logspout.conf.j2
description "logspout container"

start on started docker
stop on stopping docker

script
  /usr/bin/docker start -a logspout
end script

respawn
respawn limit 100 5
kill timeout 15
